# roles/master/kubeadm/configure/templates/etc/profile.d/kubeconfig.sh.j2

export KUBECONFIG=/etc/kubernetes/admin.conf

# vim:ts=4:sw=4:et:syn=sh:
